# crontab
15 * * * * (cd /home/finite/config; /usr/bin/ruby /home/finite/script/update_forex_data.rb)
